# etc/udev/rules.d/100-usb.rules
SUBSYSTEM=="usb", DRIVER=="usb", ACTION=="add", RUN+="/etc/init.d/mididings restart"

